# roles/mongo-config/templates/supervisor-mongod-config.conf.j2
[program:mongod-config]
command=mongod --config /etc/mongod-config.conf
user=mongodb