# roles/mongo-shard/templates/supervisor-mongod-shard.conf.j2
[program:mongod-shard]
command=mongod --config /etc/mongod-shard.conf
user=mongodb