# roles/mongo-mongos/templates/supervisor-mongos.conf.j2
[program:mongos]
command=mongos --config /etc/mongos.conf
user=mongodb